# phantom/playbooks/common/templates/phantom-maint.timer.j2
[Unit]
Description=Weekly Phantom Maintenance

[Timer]
OnCalendar=Sun *-*-* 03:00:00
RandomizedDelaySec=30m
Persistent=false

[Install]
WantedBy=timers.target
